# freehub.conf
[program:freehub]
command=/freehub/script/server
autorestart=false
